# systemd/deck.service
[Unit]
Description=FIFINE D6 Stream Deck (deck.rb)
# start after the graphical session is ready (KWin/D-Bus/Wayland up)
After=graphical-session.target
PartOf=graphical-session.target
# don't give up if the device isn't plugged/ready yet: keep retrying
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=%h/Desenvolvimento/projetos/thoths/deck/duby/bin/deck-daemon
# SIGTERM reaches ruby (cgroup) -> triggers the "Deck OFF" screen before exit
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target

# ── Install (systemd --user) ──────────────────────────────────────────────────
#   mkdir -p ~/.config/systemd/user
#   ln -sf "$PWD/systemd/deck.service" ~/.config/systemd/user/deck.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now deck.service
#   journalctl --user -u deck -f         # (in addition to deck.log in the directory)
#   systemctl --user stop deck.service   # paints "Deck OFF" and exits
#
# If kdotool/busctl/notify-send fail under the service (missing session env):
#   systemctl --user import-environment PATH WAYLAND_DISPLAY DISPLAY \
#     DBUS_SESSION_BUS_ADDRESS XDG_RUNTIME_DIR
